# apps/data-extractors/crontab
# START CRON JOB
30 12,15 * * * cd /cronjob && /usr/local/bin/python /cronjob/cronjob.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB